# examples/nordscan.service
[Unit]
Description=Nordscan

[Service]
User=nordlo
#Environment="LOG_LEVEL=DEBUG"
Environment="LOG_FILE=/var/log/nordscan/nordscan.log"
ExecStart=/opt/nordscan/.env/bin/python /opt/nordscan/src/nordscan.py --config /home/nordlo/.nordscan/config.yaml

[Install]
WantedBy=multi-user.target
